# 3.8.1/telegram-bot.service
[Unit]
Description=Telegram Gift & Sticker Bot
After=network.target
Wants=network.target

[Service]
Type=simple
User=yousefmsm1
Group=yousefmsm1
WorkingDirectory=/home/yousefmsm1/Desktop/last 3.4
Environment=PATH=/home/yousefmsm1/Desktop/last 3.4/venv/bin
ExecStart=/home/yousefmsm1/Desktop/last 3.4/venv/bin/python /home/yousefmsm1/Desktop/last 3.4/start_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-bot

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/yousefmsm1/Desktop/last 3.4

# Resource limits
LimitNOFILE=65536
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target 